# cron.d
PATH=/usr/local/sbin:/usr/sbin:/sbin:/usr/local/bin:/usr/bin:/bin                                                                        


*/2 * * * * root /usr/bin/7ea cron >/dev/null 2>&1
